# volumio/lib/systemd/system/wireless.service
[Unit]
Description=Wireless Services
[Service]
Type=forking
ExecStart=/volumio/app/plugins/system_controller/network/wireless.js start
KillMode=mixed 

[Install]
WantedBy=multi-user.target

