# volumio/buster/etc/network/if-pre-up.d/powersavedisable
#!/bin/sh
iw dev wlan0 set power_save off
exit 0
