# volumio/lib/systemd/system/netplug.service.d/volumio-network.conf
[Unit]
BindsTo=volumio-network.service

